# vendor/github.com/whosonfirst/go-whosonfirst-csvdb/init.d/wof-csvdb-server.sh.example
#!/bin/sh -e

### BEGIN INIT INFO
# Provides:          wof-csvdb-server.sh
# Required-Start:    $remote_fs $syslog
# Required-Stop:     $remote_fs $syslog
# Default-Start:     2 3 4 5
# Default-Stop:      0 1 6
# Short-Description: Start daemon at boot time
# Description:       Enable service provided by daemon.
### END INIT INFO

# update-rc.d defaults wof-csvdb-server.sh
# update-rc.d remove wof-csvdb-server.sh

# Stuff you will need to edit

CSVDB_USER=YOUR_CSVDB_SERVER_USER
CSVDB_DAEMON=YOUR_CSVDB_SERVER_PATH

CSVDB_ARGS=""

# Okay - you shouldn't need to edit anything after this

PID_PATH=/var/run/
PID_NAME=`basename ${CSVDB_DAEMON}`
PID_FILE=${PID_PATH}${PID_NAME}.pid

if [ "$(id -u)" != "0" ]; then
    echo "This script must be run as root" 1>&2
    exit 1
fi

case $1 in
    debug)

        echo "Starting ${PID_NAME} in debug mode w/ options: ${CSVDB_ARGS}"
        exec sudo -u $CSVDB_USER $CSVDB_DAEMON $CSVDB_ARGS
        ;;
    start)

    if [ -e "${PID_FILE}" ]; then
            echo "Error! ${PID_NAME} is currently running!" 1>&2
            exit 1
	    fi

        echo "Starting ${PID_NAME} w/ options: ${CSVDB_ARGS}"
        PID=`exec sudo -u $CSVDB_USER $CSVDB_DAEMON $CSVDB_ARGS  > /dev/null 2>&1 & echo $!`

        if [ -z $PID ]; then
            printf "%s\n" "Fail"
        else
            echo $PID > $PID_FILE
            printf "%s\n" "Ok"
	        echo $PID
        fi

        ;;
    stop)

    printf "%-50s" "Stopping $PID_NAME"
        PID=`cat $PID_FILE`

        if [ -f $PID_FILE ]; then
            kill -HUP $PID
            printf "%s\n" "Ok"
            rm -f $PID_FILE
        else
            printf "%s\n" "pidfile not found"
        fi

        ;;
    restart)
        $0 stop
        sleep 1
        $0 start
        ;;
    *)
        echo "Usage: $0 {start|stop|restart|debug}" >&2
        exit 1
        ;;
esac
